# deploy/bandcash.service
[Unit]
Description=BandCash web application
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/bandcash/app
ExecStart=/opt/bandcash/bin/server
Restart=always
RestartSec=3

# App environment
Environment=APP_ENV=production
Environment=PORT=2222
Environment=URL=https://bandcash.app
Environment=DB_PATH=/opt/bandcash/data/sqlite.db
Environment=LOG_LEVEL=info
Environment=LOG_FOLDER=/opt/bandcash/logs
Environment=LOG_PREFIX=bandcash

# SMTP (set real values)
Environment=SMTP_HOST=
Environment=SMTP_PORT=587
Environment=SMTP_USERNAME=
Environment=SMTP_PASSWORD=
Environment=EMAIL_FROM=BandCash <noreply@bandcash.app>

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/bandcash
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
